# systemd/presence.target
[Unit]
After=network.target
Wants=arrival.service leaving.service

[Install]
WantedBy=multi-user.target
